# bluetti-mqtt/etc/services.d/bluetti-mqtt/run
#!/usr/bin/with-contenv bash
exec python3 -m bluetti_mqtt
